# supervisord.conf
[program:uwsgi]
autostart=true
user=www
directory=%(here)s
command=%(here)s/.env/bin/uwsgi --ini uwsgi.ini
autorestart=true
stopsignal=INT
